# Sync/systemd/bandcamp-sync-reconcile.service
[Service]
Type=oneshot
ExecStart=%h/BandcampSync/bin/reconcile.sh
