# Sync/systemd/bandcamp-sync.path
[Unit]
Description=Bandcamp Sync Queue Watcher

[Path]
PathExistsGlob=%h/BandcampSync/inbox/pending/*.job

[Install]
WantedBy=default.target
