# Sync/systemd/bandcamp-sync-worker.service
[Service]
Type=oneshot
ExecStart=%h/BandcampSync/bin/worker.sh
