# Sync/systemd/bandcamp-sync-worker.path
[Unit]
Description=Watch Bandcamp Sync Queue

[Path]
PathExistsGlob=%h/BandcampSync/Sync/inbox/pending/*.job

[Install]
WantedBy=default.target
